# etc/cron.d/pfmb
0 6,12 * * * /usr/local/bin/pfmb | logger -t pfsense_backup
